# OS/Lab8/ostime.service
[Service]
ExecStart=/home/lolimpo/desktop/Projects/OS
WorkingDirectory=/home/lolimpo/desktop